# .claude/scripts/com.llmwiki.weekly-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llmwiki.weekly-digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/oracleneo/llm-wiki/.claude/scripts/weekly-digest.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/oracleneo/llm-wiki/.claude/scripts/weekly-digest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/oracleneo/llm-wiki/.claude/scripts/weekly-digest.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/oracleneo/llm-wiki</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/Users/oracleneo/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/oracleneo</string>
  </dict>
</dict>
</plist>
